# roles/observium/files/observium_agent@.service
[Unit]
Description=Observium Agent Service

[Service]
ExecStart=-/usr/bin/observium_agent
User=root
Group=root
StandardInput=socket
